# staging/systemd/nightfall-photo-ingress.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*:0/60
RandomizedDelaySec=30
Persistent=false
