# staging/systemd/nightfall-photo-ingress-api.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/env uvicorn api.app:app --host 127.0.0.1 --port 8000
