# staging/systemd/nightfall-photo-ingress.service.d/override.conf
[Unit]
Documentation=file:///opt/ingress/README.md

[Service]
ExecStart=
ExecStart=/opt/ingress/bin/nightfall-photo-ingress \
    --log-mode json poll \
    --path /etc/nightfall/photo-ingress.conf
