# metrics/systemd/nightfall-metrics-poller.service
[Unit]
Description=Nightfall Metrics Poller (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/chris/dev/nightfall-photo-ingress
ExecStart=/home/chris/dev/nightfall-photo-ingress/metricsctl run-now

[Install]
WantedBy=multi-user.target
