# staging/systemd/nightfall-photo-ingress-trash.service.d/override.conf
[Unit]
Documentation=file:///opt/ingress/README.md
After=local-fs.target
ConditionPathExists=/etc/nightfall/photo-ingress.conf

[Service]
ExecStart=
ExecStart=/opt/ingress/bin/nightfall-photo-ingress --log-mode json process-trash --path /etc/nightfall/photo-ingress.conf
WorkingDirectory=/var/lib/ingress
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ReadWritePaths=/var/lib/ingress /var/log/nightfall /var/cache/nightfall-photo-ingress /run/nightfall-status.d
MemoryMax=256M
