# systemd/nightfall-photo-ingress-trash.service
[Unit]
Description=Nightfall Photo Ingress trash processor
Documentation=https://github.com/nightfall/nightfall-photo-ingress
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/nightfall-photo-ingress/bin/nightfall-photo-ingress \
    --log-mode json process-trash \
    --path /etc/nightfall/photo-ingress.conf
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nightfall-photo-ingress-trash
TimeoutStartSec=300
MemoryMax=512M

[Install]
WantedBy=multi-user.target
